# crontab
*/15 * * * * /Users/radoslawfalak/Desktop/GIT/FrontStep/autoscript.sh
